# roles/chronograf/files/etc/systemd/system/chronograf.service.d/chronografdirs.conf
[Service]
ExecStartPre=+/usr/local/bin/chronografdirs.sh