# roles/influxdb/files/etc/systemd/system/influxdb.service.d/influxdirs.conf
[Service]
ExecStartPre=+/usr/local/bin/influxdirs.sh