# roles/influxdb/files/etc/systemd/system/influxdb.service.d/influx-schema.conf
[Service]
ExecStartPost=+/usr/local/bin/influx-schema.sh
